# app/config/_os/etc/udev/rules.d/rune_usb-audio.rules
#KERNEL=="card*", DRIVER=="snd-usb-audio", RUN+="/var/www/command/refresh_ao"
KERNEL=="card*", SUBSYSTEM=="sound", RUN+="/var/www/command/refresh_ao"
