# _OS_SETTINGS/etc/rc.local
#!/bin/sh -e

/var/www/command/player_wdog.sh startup & > /dev/null 2>&1

exit 0

